# dist/bmc-acf-manager.service
[Unit]
Description=BMC ACF manager

[Service]
ExecStart=/usr/bin/env bmc-acf
SyslogIdentifier=bmc-acf
Restart=always

Type=dbus
BusName=xyz.openbmc_project.Certs.ACF.Manager

[Install]
WantedBy=multi-user.target
